# infra/systemd/forestcatering-infra.service
[Unit]
Description=ForestCatering infrastructure stack (Docker Compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/workspace/ForestCatering/infra
ExecStart=/usr/bin/docker compose -f /workspace/ForestCatering/infra/compose.yml --env-file /workspace/ForestCatering/infra/.env up -d
ExecStop=/usr/bin/docker compose -f /workspace/ForestCatering/infra/compose.yml --env-file /workspace/ForestCatering/infra/.env down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
